# templates/init.erb
#!/bin/bash
# puppet-resque_scheduler
#
# Authors: Paul Springett
#

### BEGIN INIT INFO
# Provides:          resque_scheduler_<%= @name %>
# Short-Description: <%= @name %> running resque-scheduler
### END INIT INFO

. /etc/default/resque_scheduler_<%= @name %>

RETVAL=0

. /lib/lsb/init-functions

case "$1" in
  start)
    log_daemon_msg "Starting resque-scheduler server" "resque_scheduler_<%= @name %>"
    RAILS_ENV="<%= @rack_env %>" PIDFILE="<%= @pidfile %>" BACKGROUND=yes LOGFILE="<%= @_log_path %>" APP_NAME=<%= @name %> start-stop-daemon --start --quiet --oknodo --user "${USER:-root}" --group "${GROUP:-root}" --pidfile $PIDFILE --chdir $APPROOT --startas $DAEMON -- $DAEMON_OPTS
    RETVAL=$?
    log_end_msg $RETVAL
    ;;
  graceful-stop)
    log_daemon_msg "Stopping resque-scheduler server" "resque_scheduler_<%= @name %>"
    start-stop-daemon --stop --pidfile $PIDFILE --signal QUIT
    RETVAL=$?
    log_end_msg $RETVAL
    ;;
  quick-stop)
    log_daemon_msg "Stopping resque-scheduler server" "resque_scheduler_<%= @name %>"
    start-stop-daemon --stop --pidfile $PIDFILE --signal INT
    RETVAL=$?
    log_end_msg $RETVAL
    ;;
  stop)
    log_daemon_msg "Stopping resque-scheduler server" "resque_scheduler_<%= @name %>"
    start-stop-daemon --stop --pidfile $PIDFILE --retry=QUIT/30/INT/10/KILL/5
    RETVAL=$?
    log_end_msg $RETVAL
    ;;
  restart)
    $0 stop
    sleep 5
    $0 start
    ;;
  status)
      status_of_proc -p "$PIDFILE" 'resque_scheduler_<%= @name %>' '<%= @name %>'
      RETVAL=$?
      ;;
  *)
    echo "Usage: $0 {start|stop|graceful-stop|quick-stop|restart}"
    exit 1
    ;;
esac

exit $RETVAL
